# out_bin/v90s_thirdparty_plumos_not_by_me_20250615/userdata_system/bin/99-usb-dac.rules
ACTION=="add", SUBSYSTEM=="sound", RUN+="/userdata/system/bin/sink-reload.sh add"
ACTION=="remove", SUBSYSTEM=="sound", RUN+="/userdata/system/bin/sink-reload.sh remove"
